# supervisord_reverse_worker_example.ini
[program:reverse_worker]
command=/usr/bin/php reverse_worker.php              
process_name=%(program_name)s_%(process_num)02d
numprocs=1                    
directory=/etc/supervisord/gearman_workers/                
autostart=true                
autorestart=unexpected        
environment=GEARMAN_USER=gearmand
user=gearmand           
stdout_logfile=/var/log/supervisord/gearman/reverse_worker.log        
stopsignal=KILL
